# deploy/systemd/reconciliation-watchdog.service
[Unit]
Description=BeautyCita reconciliation watchdog — runs the 3 accounting invariants
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# CRON_SECRET must match the one in supabase-docker/.env
EnvironmentFile=/etc/beautycita/reconciliation-watchdog.env
ExecStart=/usr/bin/curl -sS --max-time 30 \
  -X POST "https://beautycita.com/supabase/functions/v1/reconciliation-watchdog" \
  -H "apikey: ${ANON_KEY}" \
  -H "Authorization: Bearer ${ANON_KEY}" \
  -H "x-cron-secret: ${CRON_SECRET}" \
  -H "Content-Type: application/json" \
  -o /var/log/reconciliation-watchdog.log
StandardOutput=journal
StandardError=journal
